# service/fcserver/run
#!/bin/sh
exec /usr/bin/fcserver /etc/fcserver.json

